# deploy/com.qtrade.paper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.qtrade.paper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/kelsey/.qtrade-paper.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>5</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/kelsey/Desktop/量化/outputs/paper/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kelsey/Desktop/量化/outputs/paper/launchd.err</string>
</dict>
</plist>
